# harbour-9store.timer
[Unit]
Description=Receive Notifications.

[Timer]
OnBootSec=60
OnUnitInactiveSec=120

[Install]
WantedBy=timers.target
