# harbour-9store.timer
[Unit]
Description=9Store timer

[Timer]
OnBootSec=120
OnUnitActiveSec=120
Unit=harbour-9store.service

[Install]
WantedBy=user-session.target
