# harbour-9store.service
[Unit]
Description=Run Notifications Job

[Service]
Type=simple
RemainAfterExit=no
ExecStart=/usr/share/harbour-9store/qml/py/jobs.sh

[Install]
WantedBy=multi-user.target
