# harbour-9store.timer
[Unit]
Description=9Store timer

[Timer]
OnBootSec=60
OnUnitInactiveSec=600
Unit=harbour-9store.service

[Install]
WantedBy=timer.target
