# harbour-9store.timer
[Unit]
Description=9Store timer

[Timer]
OnBootSec=2min
OnUnitInactiveSec=120
Unit=harbour-9store.service

[Install]
WantedBy=timer.target
