# harbour-9store.timer
[Unit]
Description=9Store timer

[Timer]
OnBootSec=120
OnUnitInactiveSec=120

[Install]
WantedBy=multi-user.target

